# deploy/wm-frontend.service
# systemd-Unit für das WM-Tippspiel-Frontend (next start).
# Ablage: /etc/systemd/system/wm-frontend.service
# Danach: systemctl daemon-reload && systemctl enable --now wm-frontend
#
# NEXT_PUBLIC_API_BASE_URL muss NICHT hier gesetzt werden — die eingecheckte
# .env.production liefert sie für Build und Laufzeit automatisch.

[Unit]
Description=WM-Tippspiel Frontend (Next.js)
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/wm-tippspiel-web
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
